# squid-game-ctf-2023/challenge-01/deploy/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:daemon]
command=/root/daemon.sh

[program:webserver]
command=/entrypoint.sh

[program:sshserver]
command=/usr/sbin/sshd -D
