# bitup-ctf-2021/fancydevs/deploy/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server

[program:apache2]
command=/usr/local/bin/apache2-foreground
